# ftp/image/scripts/jobs.crontab
0 */2 * * * rsync -rc /home/ftp_users /ftp_backups
